# src/systemd/clevis-luks-askpass.path
[Unit]
Description=Clevis systemd-ask-password Watcher

[Path]
PathChanged=/run/systemd/ask-password

[Install]
WantedBy=remote-cryptsetup-pre.target
